# packages/openwrt/etc/init.d/supernode
#!/bin/sh /etc/rc.common
START=90
STOP=10

USE_PROCD=1
PROG=/usr/bin/supernode
CONFIGFILE=/etc/n2n/supernode.conf

start_service() {
	procd_open_instance
	procd_set_param command $PROG $CONFIGFILE
	procd_set_param file $CONFIGFILE
	procd_set_param respawn
	procd_close_instance
}

stop_service()
{
	service_stop $PROG
}

service_triggers()
{
	procd_add_reload_trigger "supernode"
}
