# packages/etc/systemd/system/edge-ntopng@.service.in
[Unit]
Description=n2n edge process, on %I
After=network-online.target syslog.target
Wants=network-online.target
BindsTo=ntopng.service

[Service]
Type=simple
ExecStartPre=
ExecStart=/usr/sbin/edge /etc/n2n/edge-%i.conf -f
Restart=on-abnormal
RestartSec=5

[Install]
WantedBy=ntopng.service
Alias=
