# packages/etc/systemd/system/edge.service.in
[Unit]
Description=n2n edge process
After=network.target syslog.target
Wants=

[Service]
Type=simple
ExecStartPre=
ExecStart=/usr/sbin/edge /etc/n2n/edge.conf -f
Restart=on-abnormal
RestartSec=5

[Install]
WantedBy=multi-user.target
Alias=
