# packages/etc/systemd/system/supernode.service.in
[Unit]
Description=n2n supernode process
After=network.target syslog.target
Wants=

[Service]
Type=simple
ExecStartPre=
ExecStart=/usr/sbin/supernode /etc/n2n/supernode.conf -f
Restart=on-abnormal
RestartSec=5

[Install]
WantedBy=multi-user.target
Alias=
