# packages/debian/debian/postinst
#!/bin/sh -e

case "$1" in
    configure)
	# continue below
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
	exit 0
    ;;

    *)
	echo "postinst called with unknown argument \`$1'" >&2
	exit 0
    ;;
esac

umask 022

echo "Rebuilding ld cache..."
/sbin/ldconfig

if [ -f /.dockerenv ]; then exit 0; fi

# Start service after upgrade/install
systemctl daemon-reload
systemctl reset-failed

# Enable edge
if [ -f /etc/n2n/edge.conf ]; then
    echo "Enabling n2n edge..."
    systemctl enable edge
    systemctl restart edge
fi

# Enable supernode
if [ -f /etc/n2n/supernode.conf ]; then
    echo "Enabling n2n supernode..."
    systemctl enable supernode
    systemctl restart supernode
fi

exit 0
